# recipes-support/pir-sleep/pir-sleep/pir-sleep.service
[Unit]
Description=PIR-driven display sleep (CEC/DDC)
After=cog.service

[Service]
ExecStart=/usr/libexec/pir-sleep.sh
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
