# recipes-support/cog-watchdog/cog-watchdog/cog-watchdog.service
[Unit]
Description=Detect vc4 GPU reset cascade and force reboot

[Service]
Type=oneshot
ExecStart=/usr/libexec/cog-watchdog.sh
